# systemd/usb_modeswitch@.timer
[Unit]
Description=USB_ModeSwitch

[Timer]
OnStartupSec=5

